# rotatewallpaper.timer
[Unit]
Description=Run rotatewallpaper every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=default.target

