# rotatewallpaper.service
[Unit]
Description=Pick a new wallpaper from a configured folder

[Service]
Type=oneshot
ExecStart=/usr/bin/rotatewallpaper

